# packages/system-maintenance/package-upgrade.service
[Unit]
Description=Upgrade Packages

[Service]
Type=oneshot
ExecStart=/usr/bin/python -m "system.maintenance" "upgrade-packages"
