# packages/system-maintenance/system-maintenance@.service
[Unit]
Description=System Maintenance

[Service]
Type=oneshot
ExecStart=/usr/bin/python -m "system.maintenance" "%i"
