# packages/system-jobs/package-upgrade.service
[Unit]
Description=Upgrade Packages

[Service]
Type=oneshot
ExecStart=/usr/bin/python -m yay.jobs upgrade-packages
